# telegram-bot/systemd/telegram-agent-bot.service
# User service: Telegram bot that forwards messages to Cursor agent.
# Install: copy to ~/.config/systemd/user/
#   cp telegram-bot/systemd/telegram-agent-bot.service ~/.config/systemd/user/
# Then edit the paths below if your repo is not in ~/projects/cursor-claw.
# Enable lingering so it runs without a login session:
#   loginctl enable-linger $USER
# Start:
#   systemctl --user enable --now telegram-agent-bot.service

[Unit]
Description=Telegram agent bot (Cursor)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/cursor-claw
ExecStart=/usr/bin/python3 %h/projects/cursor-claw/telegram-bot/agent_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
